# automation/launchd/com.allexishere.quality.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.allexishere.quality</string>
  <key>ProgramArguments</key><array><string>/opt/homebrew/bin/node</string><string>/Users/jeonghunoh/Projects/allexishere-blog/automation/quality-review.mjs</string></array>
  <key>WorkingDirectory</key><string>/Users/jeonghunoh/Projects/allexishere-blog</string>
  <!-- 매주 목요일 09:00. 품질 검토는 '기존 글이 괜찮은가'(유지) — 주간 1회.
       월요일 SEO 감시와 다른 요일에 둬 부하·주목을 분산. 조치는 승인 후. -->
  <key>StartCalendarInterval</key><array>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>EnvironmentVariables</key><dict><key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin</string></dict>
  <key>StandardOutPath</key><string>/Users/jeonghunoh/Projects/allexishere-blog/automation/logs/quality.out.log</string>
  <key>StandardErrorPath</key><string>/Users/jeonghunoh/Projects/allexishere-blog/automation/logs/quality.err.log</string>
</dict></plist>
